# rootfs/var/lib/dpkg/info/cgroupfs-mount.list
/.
/etc
/etc/init
/etc/init/cgroupfs-mount.conf
/etc/init.d
/etc/init.d/cgroupfs-mount
/usr
/usr/share
/usr/share/doc
/usr/share/doc/cgroupfs-mount
/usr/share/doc/cgroupfs-mount/README.md
/usr/share/doc/cgroupfs-mount/copyright
/usr/share/doc/cgroupfs-mount/changelog.gz
/usr/bin
/usr/bin/cgroupfs-umount
/usr/bin/cgroupfs-mount
