# rootfs/var/lib/dpkg/info/cgroup-lite.list
/.
/etc
/etc/init
/etc/init/cgroup-lite.conf
/bin
/bin/cgroups-umount
/bin/cgroups-mount
/lib
/lib/systemd
/lib/systemd/system
/usr
/usr/bin
/usr/share
/usr/share/doc
/usr/share/doc/cgroup-lite
/usr/share/doc/cgroup-lite/changelog.gz
/usr/share/doc/cgroup-lite/copyright
/lib/systemd/system/cgroup-lite.service
/usr/bin/cgroups-umount
/usr/bin/cgroups-mount
